# ansible/templates/mystery-music.service.j2
[Unit]
Description={{ service_description | default('Mystery Melody Machine') }}
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User={{ pi_user }}
Group=audio
WorkingDirectory={{ rpi_dir }}
Environment=PATH={{ venv_dir }}/bin
ExecStartPre=/bin/bash -c 'modprobe snd-seq || true'
ExecStart={{ rpi_dir }}/start-mystery-music.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10
StandardOutput=append:/var/log/mystery-melody/mystery-melody.log
StandardError=append:/var/log/mystery-melody/mystery-melody-error.log

# Security settings (relaxed for GPIO access on Zynthian)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ rpi_dir }} /var/log/mystery-melody {% if is_zynthian_hardware %}/dev/gpiomem /sys/class/gpio{% endif %}

# Group permissions for hardware access
SupplementaryGroups=audio,dialout,plugdev{% if is_zynthian_hardware %},gpio,i2c,spi{% endif %}

# Resource limits (increased for Zynthian performance)
LimitNOFILE=2048
MemoryMax=1G
{% if is_zynthian_hardware %}
# Real-time scheduling for low-latency hardware response
LimitRTPRIO=95
LimitMEMLOCK=infinity
{% endif %}

[Install]
WantedBy=multi-user.target
